# conf/sabnzbd.ini
[program:sabnzbd]
command = /opt/sabnzbd/SABnzbd.py -l0 -f /srv/sabnzbd/sabnzbd.ini
directory = /opt/sabnzbd
user = sabnzbd
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/sabnzbd.log
stderr_logfile = /var/log/supervisor/sabnzbd_errors.log
environment =
