# conf/sickrage.ini
[program:sickrage]
command = /usr/bin/env python2 /opt/sickrage/SickBeard.py --quiet --config /srv/sickrage/config.ini --datadir /srv/sickrage
directory = /opt/sickrage
user = sickrage
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/sickrage.log
stderr_logfile = /var/log/supervisor/sickrage_errors.log
environment =
